# scripts/crontab
55 23 * * *     /var/lib/postgresql/exfetch/scripts/ext_table_and_restart.sh >> /tmp/e.log 2>&1
*/10 * * * *       /var/lib/postgresql/exfetch/scripts/bs_in_a_row.sh >> /tmp/f.log 2>&1
